# package/upapps/luci-theme-oceanic/root/etc/uci-defaults/30_luci-theme-material
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Material=/luci-static/oceanic
	set luci.main.mediaurlbase=/luci-static/oceanic
	commit luci
EOF
exit 0
